# deploy/templates/hub.service
[Unit]
Description= Hub - Alliance Hub for pr0game
After=network.target

[Service]
Type=simple
User=hub
Group=hub
WorkingDirectory=/opt/hub
ExecStart=/opt/hub/bin/hub
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hub

# Environment
EnvironmentFile=/opt/hub/.env

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/hub/data /opt/hub/logs

[Install]
WantedBy=multi-user.target
